# ops/systemd/vpn-bot-web.service
[Unit]
Description=VPN Telegram Bot Webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vpn-bot
Environment=VPN_BOT_HOME=/opt/vpn-bot
Environment=VPN_BOT_SECRETS_FILE=/opt/vpn-bot/secrets/runtime.toml
ExecStart=/opt/vpn-bot/.venv/bin/python -m vpn_bot web
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
